# etc-udev/1-drinkallthemilk.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/usr/sbin/datm-ntfy"
ACTION=="remove", SUBSYSTEM=="usb", RUN+="/usr/sbin/datm-ntfy"
